# deploy/celery/supervisord.conf
; Celery worker + beat supervision for bare-metal / VM deployments
; Usage: supervisord -c deploy/celery/supervisord.conf

[supervisord]
nodaemon=true
logfile=/var/log/medsync/supervisord.log
pidfile=/var/run/medsync/supervisord.pid
childlogdir=/var/log/medsync

[program:celery-worker]
command=celery -A medsync_backend worker --loglevel=info --concurrency=4
directory=/app/medsync-backend
user=medsync
autostart=true
autorestart=true
stopwaitsecs=600
stdout_logfile=/var/log/medsync/celery-worker.log
stderr_logfile=/var/log/medsync/celery-worker.err.log
environment=PYTHONUNBUFFERED="1"

[program:celery-beat]
command=celery -A medsync_backend beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/app/medsync-backend
user=medsync
autostart=true
autorestart=true
stdout_logfile=/var/log/medsync/celery-beat.log
stderr_logfile=/var/log/medsync/celery-beat.err.log
environment=PYTHONUNBUFFERED="1"
